# ops/ansible_project/roles/webserver/templates/app.conf.j2
[program:{{ app_name }}]
command = gunicorn --bind {{ server_ip }}:{{ server_port }} {{ app_name }}:app
directory = {{ app_dir }}